# deploy/systemd/cost-update.service
[Unit]
Description=Materialize daily energy kWh and compute electricity cost metrics
After=victoria-metrics.service

[Service]
Type=oneshot
WorkingDirectory=/opt/homeserver
ExecStart=/home/homepi/.local/bin/uv run el cost-update
User=homepi
Group=homepi
