# deploy/systemd/rpi-metrics.service
[Unit]
Description=Export Raspberry Pi throttling metrics for node_exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rpi-metrics.sh
