# deploy/systemd/switchbot.service
[Unit]
Description=SwitchBot sensor collector
After=network-online.target bluetooth.service victoria-metrics.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
WorkingDirectory=/opt/homeserver
ExecStart=/home/homepi/.local/bin/uv run sb run --interval 60 --mode ble --ble-scan-timeout 20
Restart=on-failure
RestartSec=30
User=homepi
Group=homepi

[Install]
WantedBy=multi-user.target
