# deploy/systemd/vm-backup.service
[Unit]
Description=Backup VictoriaMetrics snapshot to pCloud
After=victoria-metrics.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vm-backup.sh
