# deploy/collector-watchdog.timer
[Unit]
Description=Sensor data freshness watchdog

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
